# supervisord-syslog_to_s3.conf
[program:syslog_to_s3]
command=/start-syslog_to_s3.sh
numprocs=1
autostart=true
autorestart=true
